# StockProjectBackend/crontab
*/10 18-19 * * * /usr/local/bin/python /app/cronjobs/run_cron_task.py >> /var/log/cron.log 2>&1
